# supervisor.conf
;
; supervisor.conf — Mini CRM Queue Worker
;
; Copy this file to /etc/supervisor/conf.d/mini-crm.conf
; then run: supervisorctl reread && supervisorctl update && supervisorctl start mini-crm-worker:*
;
; This keeps the Laravel queue worker running in the background so that
; queued notifications (email assignments, task reminders) are delivered.
;
; Author  : Xgenious (https://xgenious.com)
; License : MIT
;

[program:mini-crm-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work database --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/worker.log
stopwaitsecs=3600

[program:mini-crm-scheduler]
process_name=%(program_name)s_%(process_num)02d
; The scheduler tick runs every minute. Laravel checks whether any scheduled
; commands are due and executes them (e.g. task due reminders at 08:00).
command=php /var/www/html/artisan schedule:work
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/scheduler.log
